# pdns-admin/supervisor.ini
[supervisord]
nodaemon = true
user = root
logfile = /dev/stderr
logfile_maxbytes = 0
loglevel = info

[program:uwsgi]
command = /usr/sbin/uwsgi --ini /etc/uwsgi.ini
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:caddy]
command = /usr/bin/caddy run --config /etc/caddy/Caddyfile
user = caddy
environment = HOME="/var/lib/caddy"
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
